# extras/whisper-dictation/config/systemd/user/ydotool.service
[Unit]
Description=ydotoold - ydotool daemon

[Service]
Type=simple
Restart=always
ExecStart=/usr/bin/ydotoold

[Install]
WantedBy=default.target
